# files/etc/supervisor.d/mysql.ini
[program:mysql]
command=/usr/bin/pidproxy /data/mysql/run.pid /usr/bin/mysqld --basedir=/usr --datadir=/data/mysql --plugin-dir=/usr/lib/mariadb/plugin --user=mysql --log-error=/data/mysql/log.err --pid-file=run.pid %CMD%
autostart=true
autorestart=true
user=root
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0